# docker/runit/monit/run
#!/bin/sh

erb /app/docker/monitrc.erb > /app/docker/monitrc

chmod 0700 /app/docker/monitrc

exec /usr/bin/monit -I -c /app/docker/monitrc
