# docker/runit/monit/run
#!/bin/sh

exec /usr/bin/monit -I -c /app/docker/monitrc
