# upstart/delayed_job.conf
setuid memecaptain

start on runlevel [2345]
stop on runlevel [016]
respawn

chdir /var/www/memecaptain.com

script
  /bin/bash -l -c 'RAILS_ENV=production script/delayed_job run'
end script
